# crontab
*      *       *       *       *       PATH=/usr/local/bin:/usr/local/sbin:~/bin:/usr/bin:/bin:/usr/sbin:/sbin && /Users/tschmelzer/cpu_monitoring/cpu_analysis.py
